# waveshare.service
[Unit]
Description=Waveshare 7inch userspace driver

[Service]
Type=simple
ExecStart/usr/bin/waveshare

[Install]
WantedBy=multi-user.target
